# rootfs/etc/s6-overlay/s6-rc.d/mkfifo/up
#!/command/execlineb -P

rm -f /tmp/audio.fifo && mkfifo /tmp/audio.fifo
